# bin/statswatch.service
[Unit]
Description=Watch service for Zend Framework web site statistics

[Service]
ExecStart=/var/vhosts/framework.zend.com/htdocs/bin/watch.js
Restart=always
User=www-data
Group=www-data
Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
WorkingDirectory=/var/vhosts/framework.zend.com/htdocs

[Install]
WantedBy=multi-user.target
