# .aws/www-data.crontab
# m h day-of-month month day-of-week command
0 0 * * 5 /var/vhosts/framework.zend.com/htdocs/bin/zfweb.php lts:build /var/vhosts/framework.zend.com/htdocs/data/long-term-support.json
