# .deploy/www-data.crontab
HOME=/var/vhosts
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin:/snap/bin

# m h day-of-month month day-of-week command
0 0 * * 5 cd /var/vhosts/framework.zend.com/htdocs ; php bin/zfweb.php lts:build
